# deploy/debug/systemd/user/invoiceinator-vue.service
[Unit]
Description=Invoiceinator Vue debug server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Dev/Invoiceinator/invoice-frontend
ExecStartPre=/usr/bin/bash -lc 'source "$HOME/.nvm/nvm.sh" && npm ci'
ExecStart=/usr/bin/bash -lc 'source "$HOME/.nvm/nvm.sh" && npm run dev'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
