# deploy/debug/systemd/user/invoiceinator-django.service
[Unit]
Description=Invoiceinator Django debug server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Dev/Invoiceinator/invoiceinator
Environment=PYENV_ROOT=%h/.pyenv
Environment=PYENV_VERSION=invoiceinator
Environment=PATH=%h/.pyenv/bin:%h/.pyenv/shims:/usr/local/bin:/usr/bin:/bin
Environment=PYTHONUNBUFFERED=1
ExecStart=%h/.pyenv/bin/pyenv exec python manage.py runserver 0.0.0.0:9999
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
